# k3s/bootstrap/systemd/network-watchdog.service
[Unit]
Description=Network connectivity watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/network-watchdog.env
ExecStart=/usr/local/bin/network-watchdog.sh
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
StartLimitIntervalSec=0
